# php-fpm/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command = /usr/sbin/php5-fpm -FR
user = root
autostart = true
